# deploy/compose/podman/cdc-ats-nats.container
# Quadlet unit - NATS JetStream (durable pipeline events).
[Unit]
Description=CDC ATS NATS JetStream
After=network-online.target

[Container]
ContainerName=nats
Image=docker.io/library/nats:2.10-alpine
Network=cdc-ats.network
Volume=cdc-ats-natsdata:/data
Exec=-js -sd /data -m 8222

[Service]
Restart=unless-stopped

[Install]
WantedBy=default.target
